# packages/niri/files/config/systemd/user/niri.service.d/override.conf
[Service]
# Keep Niri session env owned by the Niri package, not the display manager.
# Import before `niri --session` so niri, niri-spawned apps, user services,
# and D-Bus activation all see the same static session vars.
ExecStartPre=sh -c 'script="$1"; env_file="$2"; [ ! -f "$script" ] || exec sh "$script" "$env_file"' _ %h/.config/niri/bin/import-session-environment.sh %h/.config/niri/session.env
{% if vars.niri.exec_start %}
ExecStart=
ExecStart={{ vars.niri.exec_start }}
{% endif %}
