# packages/xsettings/files/config/systemd/user/xsettingsd.service.d/session.conf
[Unit]
After=graphical-session.target

[Service]
# Xwayland is spawned outside systemd by xwayland-satellite, so session.target
# does not guarantee DISPLAY is connectable yet. Retry instead of encoding fake
# ordering against a non-unit.
Restart=on-failure
RestartSec=1
