# packages/linux/desktop/files/config/autostart/goldendict.desktop
[Desktop Entry]
Categories=Office;Dictionary;Education;Qt;
Comment=A feature-rich dictionary lookup program
Exec=goldendict
GenericName=Multiformat Dictionary
Icon=goldendict
MimeType=x-scheme-handler/goldendict;x-scheme-handler/dict;
Name=GoldenDict-ng
StartupWMClass=GoldenDict-ng
Terminal=false
Type=Application
X-Desktop-File-Install-Version=0.27
