# packages/wayvnc/files/config/systemd/user/wayvnc-novnc.service
[Unit]
Description=noVNC browser bridge for WayVNC
After=wayvnc.service
BindsTo=wayvnc.service
PartOf=wayvnc.service

[Service]
Type=simple
ExecStart=/usr/bin/novnc --listen 127.0.0.1:6080 --vnc 127.0.0.1:5900
Restart=on-failure
RestartSec=3
TimeoutStopSec=10

[Install]
# Bridge is enabled under wayvnc.service, not graphical-session.target, to avoid
# a start transaction cycle with wayvnc.service's After=graphical-session.target.
WantedBy=wayvnc.service
