# packages/keyd/files/config/systemd/user/keyd-application-mapper.service
[Unit]
Description=keyd application mapper
PartOf=graphical-session.target
After=graphical-session.target
Requisite=graphical-session.target
StartLimitIntervalSec=300
StartLimitBurst=3

[Service]
ExecStart=keyd-application-mapper
Restart=on-failure
RestartSec=5

[Install]
WantedBy=graphical-session.target
