# dotfiles/config/systemd/user/sunshine.service.d/override.conf
[Unit]
After=graphical-session.target
; After=xdg-desktop-autostart.target
[Service]
StandardOutput=null
StandardError=null
[Install]
; WantedBy=graphical-selssion.target
WatedBy=default.target
