# packages/niri/files/config/systemd/user/niri-shell.service
[Unit]
Description=Noctalia shell for Niri
After=niri.service
BindsTo=niri.service
PartOf=niri.service
Before=xdg-desktop-autostart.target
Before=graphical-session.target
StartLimitIntervalSec=5min
StartLimitBurst=3

[Service]
Type=exec
TimeoutStartSec=20
ExecStart=/usr/bin/qs -c noctalia-shell
ExecStartPost=/usr/local/bin/wait-for-status-notifier-host
Restart=on-failure

[Install]
WantedBy=niri.service

