# dotfiles/config/systemd/user/keyd-application-mapper.service
[Unit]
Description=keyd application mapper
[Service]
ExecStart=keyd-application-mapper
Restart=always
RestartSec=5
[Install]
WantedBy=default.target
