# packages/xsettings/files/config/systemd/user/xsettingsd-session.target
[Unit]
Description=Pull xsettingsd into the graphical session
Wants=xsettingsd.service
PartOf=graphical-session.target

[Install]
WantedBy=graphical-session.target
